# mysql/mysqld.service
[Unit]
Description=MySQL Server

[Service]
User=root
ExecStart=/usr/sbin/mysqld --user=root
Restart=always

[Install]
WantedBy=multi-user.target
